# systemd/.config/systemd/user/hyprland.service
[Unit]
Description=The best tiling window manager for wayland

[Service]
ExecStart=Hyprland

[Install]
WantedBy=graphical-session.target
