# systemd/.config/systemd/user/nm-applet.service
[Unit]
Description=NetworkManager applet
PartOf=hyprland-session.target

[Service]
ExecStart=nm-applet

[Install]
WantedBy=hyprland-session.target

