# systemd/.config/systemd/user/waybar.service
[Unit]
Description=Title bar
PartOf=hyprland-session.target

[Service]
Type=simple
ExecStart=waybar

[Install]
WantedBy=hyprland-session.target
